# crontab
# START CRON JOB
* 1 * * * python3 /app/FetchNetflixTop10.py > /proc/1/fd/1 2>/proc/1/fd/2
# END CRON JOB